# services/next_photo.timer
[Unit]
Description=Run Next Photo Service daily at 5 AM

[Timer]
OnCalendar=*-*-* 05:00:00
Unit=next_photo.service

[Install]
WantedBy=timers.target