# services/wifi-power-save-off.service
[Unit]
Description=Disable WiFi power saving

[Service]
Type=oneshot
ExecStart=/sbin/iw dev wlan0 set power_save off
